# sys-fs/sanoid/files/sanoid.timer
[Unit]
Description=sanoid process timer
Documentation=man:sanoid(1)

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
